# startup/startup.service
[Unit]
Description=CFSD Start script

[Service]
ExecStart=/home/cfsd/boot/startup.sh

[Install]
WantedBy=multi-user.target 